# etc/ipmidetectd.service.in
[Unit]
Description=IPMI Node Detection Monitoring Daemon
After=network.target

[Service]
Type=forking
ExecStart=@X_SBINDIR@/ipmidetectd

[Install]
WantedBy=multi-user.target
